# dot_config/systemd/user/battery-log.service
[Service]
ExecStart=%h/.config/scripts/battery-log.sh
